# crontab/loona.sh
30 00 * * * sh /home/loona/programming/scripts/journal_cleaner.sh
